# sandbox/s6/init-state-perms/up
foreground {
    mkdir -p /state/worker /state/daemon
}
chown -R alice:alice /state
