# sandbox/worker/s6/alice-speaking/run
#!/command/with-contenv bash
# Speaking Alice — Python daemon (Claude Agent SDK) that replaces the legacy
# signal-bridge.sh. Runs as alice.
#
# Holds the flock on /state/worker/lease for the process lifetime so the
# existing blue/green deploy mechanism (alice-deploy worker) still coordinates
# handoff correctly. When blue/green is retired the flock wrapper can go.
set -e

# Skip the signal-cli readiness probe when Signal is disabled — the daemon
# container is up but signal-cli isn't running, so polling would just burn
# 60s before timing out on every (re)start.
SIGNAL_DISABLED=
if [ -f /state/daemon/signal-disabled ]; then
    SIGNAL_DISABLED=1
fi

if [ -z "${SIGNAL_DISABLED}" ]; then
    # Wait for signal-cli's HTTP endpoint. depends_on's service_healthy normally
    # handles this but belt-and-braces.
    for _ in $(seq 1 60); do
        if curl -sfo /dev/null -X POST "${SIGNAL_API:-http://alice-daemon:8080}/api/v1/rpc" \
            -H 'Content-Type: application/json' \
            -d '{"jsonrpc":"2.0","method":"version","id":"ping"}' 2>/dev/null; then
            break
        fi
        sleep 1
    done
fi

export ALICE_CONFIG=/home/alice/.config/alice/alice.env
# Monorepo root bind-mounted at /home/alice/alice — src/ holds all Python
# packages (alice_core, alice_speaking, alice_thinking, alice_viewer).
export PYTHONPATH=/home/alice/alice/src
# Unbuffered stdout so s6 can see logs in real time.
export PYTHONUNBUFFERED=1

# Redirect stderr to a host-visible log file so install/runtime issues are
# debuggable from `tail -F ~/.local/state/alice/worker/speaking-stderr.log`.
# stdlib logging defaults to stderr; the EventLogger writes to /state/speaking.log.
#
# s6-redirfd opens fd 2 *after* s6-setuidgid drops privileges, so the log
# file is created/owned by alice. A shell-level `2>>` would be evaluated by
# the parent shell (still root), leaving the inode root-owned and breaking
# subsequent alice-uid writes (#10).
#
# Hot-restart loop: the inner `while true` re-execs the Python interpreter
# when /state/worker/reload-requested-internal exists. The alice-reload-watcher
# s6 service touches that flag and then runs `s6-svc -r` against this
# service; the SIGTERM exits the Python process, the loop sees the flag,
# clears it, and starts a fresh interpreter — picking up bind-mounted code
# changes without dropping the flock or recycling the s6 service. If the
# flag isn't set after the Python process exits, this is a normal shutdown
# (or a crash) and we break out so s6's own restart-backoff handles it.
# Design: cortex-memory/research/2026-05-05-speaking-self-deploy-design.md.
exec s6-setuidgid alice s6-redirfd -a 2 /state/worker/speaking-stderr.log \
    flock -x /state/worker/lease \
    bash -c '
        while true; do
            /opt/alice-venv/bin/python -m alice_speaking
            if [[ -f /state/worker/reload-requested-internal ]]; then
                rm -f /state/worker/reload-requested-internal
            else
                break
            fi
        done
    '
