# sandbox/daemon/s6/init-state-perms/up
foreground {
    mkdir -p /state/daemon
}
chown alice:alice /state/daemon
