# sandbox/host-claude-watcher/alice-host-claude.service
[Unit]
Description=Alice host-side Claude task watcher (Speaking <-> Claude CLI bridge)
Documentation=https://github.com/jcronq/alice/tree/master/sandbox/host-claude-watcher
After=network.target

[Service]
Type=simple
# Adjust User= to the account that owns ~/.claude and has the `claude`
# binary on PATH. Don't run this as root.
User=alice
Group=alice
# ExecStart runs the script directly from the bind-mounted repo so
# `git pull && systemctl restart alice-host-claude` is the canonical
# update path — no copy-to-/usr/local/bin step. Issue #144: prior layout
# (`ExecStart=/usr/local/bin/alice-host-claude-watcher.sh`) silently
# left the deployed script stale after every repo change, so PR #138's
# inotify fix never took effect on the host.
#
# %h is a systemd specifier that expands to User='s home dir at unit
# load time. We use it directly in ExecStart because the first arg of
# ExecStart must be an absolute path or simple filename — env-var
# expansion (`${ALICE_REPO}/...`) is *rejected* by systemd as "Neither
# a valid executable name nor an absolute path", which silently left
# the service un-restartable after PR #147. If your clone lives
# somewhere other than ~/alice, override ExecStart= via a drop-in.
# The script discovers ALICE_HOST_CLAUDE_ROOT (default: /state/worker/host-claude).
# Override here if your bind-mount layout differs.
# Environment=ALICE_HOST_CLAUDE_ROOT=/state/worker/host-claude
# Environment=ALICE_HOST_CLAUDE_BIN=/usr/local/bin/claude
ExecStart=%h/alice/sandbox/host-claude-watcher/alice-host-claude-watcher.sh
Restart=on-failure
RestartSec=5
# Capture all log lines to the journal; the script logs to stderr.
StandardOutput=journal
StandardError=journal
# Hardening — tasks only ever touch ALICE_HOST_CLAUDE_ROOT and PATH-located
# binaries. Tighten further once the deployment shape settles.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/state/worker/host-claude
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
