# sandbox/worker/s6/alice-reload-watcher/run
#!/command/with-contenv bash
# alice-reload-watcher — watches /state/worker/reload-requested and triggers
# a hot restart of the alice-speaking Python process via the in-script
# restart loop (see sandbox/worker/s6/alice-speaking/run).
#
# Hot restart only recycles the Python interpreter inside the existing
# container — bind-mounted source code is already present, only stale
# imports need to clear. Container-level deploys remain manual.
#
# Trigger contract:
#   - Speaking writes JSON to /state/worker/reload-requested via the
#     mcp__alice__request_worker_reload tool.
#   - Watcher fires on inotify create/modify, sleeps 3s grace, then touches
#     /state/worker/reload-requested-internal which causes the alice-speaking
#     run script's `while true` loop to re-exec the Python interpreter after
#     SIGTERM.
#   - Rate limit: drop reload requests <60s after the last accepted one.
#     Prevents bad-startup loops.
#
# Design: cortex-memory/research/2026-05-05-speaking-self-deploy-design.md
# Gap fixes: cortex-memory/research/2026-05-05-speaking-self-deploy-impl-gaps.md
#
# Runs as alice — the state volume is alice-owned and inotifywait does not
# need root.

set -e

SENTINEL=/state/worker/reload-requested
INTERNAL_FLAG=/state/worker/reload-requested-internal
LAST_TS=/state/worker/.last-reload-ts
RATE_LIMIT_SEC=60
GRACE_SEC=3

# Startup drain: if a sentinel exists from a previous run (e.g. watcher
# crashed before it could process), clean it up and act on it once.
# inotifywait does NOT fire for files that already exist when it starts.
if [[ -f "$SENTINEL" ]]; then
    echo "alice-reload-watcher: stale sentinel found at startup; processing once." >&2
    rm -f "$SENTINEL"
    sleep "$GRACE_SEC"
    touch "$INTERNAL_FLAG"
    s6-svc -t /run/service/alice-speaking
fi

# Ensure /state/worker exists; init-state-perms is a oneshot we depend on,
# but be defensive if anything reorders.
mkdir -p /state/worker

exec s6-setuidgid alice \
    inotifywait -e create -e modify -m /state/worker/ --format '%f' \
    | while read -r fname; do
        [[ "$fname" = "reload-requested" ]] || continue

        # Read payload (best-effort); empty/invalid is OK, we still proceed.
        payload=$(cat "$SENTINEL" 2>/dev/null || echo '{}')
        rm -f "$SENTINEL"

        # Rate limit: if the last accepted reload was <RATE_LIMIT_SEC ago,
        # drop the request. Guards against bad-code restart loops.
        now_epoch=$(date +%s)
        if [[ -f "$LAST_TS" ]]; then
            last_epoch=$(cat "$LAST_TS" 2>/dev/null || echo 0)
            since_last=$((now_epoch - last_epoch))
            if (( since_last < RATE_LIMIT_SEC )); then
                echo "alice-reload-watcher: rate-limited (${since_last}s since last reload < ${RATE_LIMIT_SEC}s); dropping." >&2
                continue
            fi
        fi
        echo "$now_epoch" > "$LAST_TS"

        # Reload type — only "hot" is supported in v1. Default to hot.
        type=$(printf '%s' "$payload" | python3 -c "import sys,json; sys.stdout.write(str(json.load(sys.stdin).get('type','hot')))" 2>/dev/null || echo "hot")
        if [[ "$type" != "hot" ]]; then
            echo "alice-reload-watcher: unsupported reload type '$type'; dropping." >&2
            continue
        fi

        echo "alice-reload-watcher: hot-reload requested; sleeping ${GRACE_SEC}s for grace period." >&2
        sleep "$GRACE_SEC"

        # Touch the internal flag so the alice-speaking restart loop knows
        # to re-exec rather than break out. Then SIGTERM the Python process
        # via s6-svc -t; the run script's bash-c while-true catches the
        # internal flag and starts a new interpreter.
        #
        # (Originally used `s6-svc -r`; verified empirically 2026-05-26 on
        # s6 2.14 that -r returned rc=0 but did NOT terminate the child.
        # -t sends SIGTERM as the run-script's trap expects. Same fix as
        # #392 for the cozylobe-reload-watcher.)
        touch "$INTERNAL_FLAG"
        s6-svc -t /run/service/alice-speaking
        echo "alice-reload-watcher: SIGTERM sent to alice-speaking." >&2
    done
