# sandbox/worker/s6/init-state-perms/up
foreground {
    mkdir -p /state/worker
}
chown alice:alice /state/worker
