# sandbox/s6/alice-cozylobe/run
#!/command/with-contenv bash
# alice-cozylobe — CozyHem reasoning lobe daemon. Long-lived Python
# process that opens an SSE connection to cozyhem-engine, runs a
# wake loop against the resulting events, and emits structured
# events to /state/worker/cozylobe.log via core.EventLogger.
#
# Runs as alice. Single supervised process — no worker-lease or
# blue/green coordination (it's not part of the speaking/thinking
# turn-handling path). Crash semantics: any task dying triggers an
# orderly shutdown; s6's own restart backoff brings it back up.
#
# Service scaffold deferred from PR #359 (cozylobe walking skeleton);
# wired here per the daemon docstring's "follow-up PR" note.
#
# Hot-restart loop (#390): the inner `while true` re-execs the Python
# interpreter when /state/worker/cozylobe-reload-requested-internal exists.
# The alice-cozylobe-reload-watcher s6 service touches that flag and
# then runs `s6-svc -r` against this service; SIGTERM forwards to the
# Python process, the loop sees the flag, clears it, and starts a fresh
# interpreter — picking up bind-mounted code changes without recycling
# the s6 service. If the flag isn't set after the Python process exits,
# this is a normal shutdown (or a crash) and we break out so s6's own
# restart-backoff handles it. Mirrors the speaking pattern documented in
# cortex-memory/research/2026-05-05-speaking-self-deploy-design.md.
set -e

export ALICE_CONFIG=/home/alice/.config/alice/alice.env
# Monorepo root bind-mounted at /home/alice/alice — src/ holds all Python
# packages (core, alice_speaking, alice_thinking, alice_cozylobe).
export PYTHONPATH=/home/alice/alice/src
# Unbuffered stdout so s6 can see logs in real time.
export PYTHONUNBUFFERED=1

# redirfd (execline) opens fd 2 *after* s6-setuidgid drops privileges so
# the stderr log is created/owned by alice. A shell-level `2>>` would be
# opened by the supervising root shell and leave the inode root-owned,
# breaking subsequent alice-uid writes. Mirrors the pattern in
# alice-speaking/run and alice-sm-dispatcher/run.
exec s6-setuidgid alice redirfd -a 2 /state/worker/cozylobe-stderr.log \
    bash -c '
        shutdown=
        child=
        on_signal() {
            shutdown=1
            if [[ -n "$child" ]]; then
                kill -TERM "$child" 2>/dev/null || true
            fi
        }
        trap on_signal TERM INT

        while true; do
            /opt/alice-venv/bin/python -m alice_cozylobe &
            child=$!

            # Signals interrupt wait (return 128+signum); the trap then
            # forwards SIGTERM to Python. We re-wait until Python is
            # actually reaped — drain may take a moment. ``|| rc=$?``
            # keeps the non-zero from killing the script under set -e.
            while :; do
                rc=0
                wait "$child" 2>/dev/null || rc=$?
                kill -0 "$child" 2>/dev/null || break
            done

            if [[ -n "$shutdown" ]]; then
                # Real shutdown: do not honour a stale reload flag.
                break
            fi
            if [[ -f /state/worker/cozylobe-reload-requested-internal ]]; then
                rm -f /state/worker/cozylobe-reload-requested-internal
            else
                break
            fi
        done
        exit "${rc:-0}"
    '
